# doc/supervisor/conf.d/zetcd.conf
[program:zetcd]
command=/usr/etcd/zetcd --zkaddr 0.0.0.0:2181 --endpoints 127.0.0.1:2379 ; the program (relative uses PATH, can take args)
;directory=/tmp                ; directory to cwd to before exec (def no cwd)
;umask=022                     ; umask for process (default None)
;priority=999                  ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
startsecs=3                    ; # of secs prog must stay up to be running (def. 1)
;startretries=3                ; max # of serial start failures when starting (default 3)
;autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
;exitcodes=0,2                 ; 'expected' exit codes used with autorestart (default 0,2)
stopsignal=INT                 ; signal used to kill process (default TERM)
;stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
;redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/dev/null       ; stdout log path, NONE for none; default AUTO
stderr_logfile=/dev/null       ; stderr log path, NONE for none; default AUTO
environment=GOTRACEBACK="crash" ; process environment additions (def no adds)
